# deploy/withssl.service
[Unit]
Description=WithSSL SSL Proxy Server
After=network.target

[Service]
Type=simple
User=withssl
Group=withssl
WorkingDirectory=/opt/withssl
ExecStart=/opt/withssl/withssl --config /etc/withssl/withssl.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=withssl

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/withssl /etc/withssl /opt/withssl

[Install]
WantedBy=multi-user.target
